# build/wg-vpn/etc/services.d/knot-resolver/run
#!/command/with-contenv bash

### KNOT
#if [[ ! -d "/storage/confdb" ]]; then
#    knotc conf-init
#fi

if [[ ! -f "/etc/knot-resolver/aliases/aliases.conf" ]]; then
    cp /knot_aliases_default.conf /etc/knot-resolver/aliases/aliases.conf
fi
# Remove knot tmp and runtime
rm -rf /tmp/knot-resolver/*
rm -rf /run/knot-resolver/*

iptables -t nat -N dnsmap
iptables -t nat -A PREROUTING -s 10.224.0.0/15 -d 10.224.0.0/15 -j dnsmap

kresd -c /etc/knot-resolver/kresd.conf -n /run/knot-resolver
