# build/wg-vpn/etc/services.d/dnsmap/run
#!/command/with-contenv bash
cd /opt/dnsmap && nohup ./proxy.py -a 127.0.0.4 --iprange 10.224.0.0/15
